# rootfs/etc/s6-overlay/s6-rc.d/ninja/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Invoice Ninja
# Cleanup and graceful shutdown
# ==============================================================================

bashio::log.info "Invoice Ninja is stopping (finish script)..."

# Der php artisan serve-Prozess sollte bereits durch s6 beendet worden sein.
# Das explizite pkill wird hier entfernt, da es unnötig ist.

# Cache und temporäre Dateien bereinigen
cd /var/www/app || bashio::exit.ok # Wechsel ins Arbeitsverzeichnis
if [ -d storage/framework/cache ]; then
    bashio::log.info "Clearing application cache..."
    # Verwende rm -r * anstelle von rm -rf, um das Verzeichnis selbst nicht zu löschen
    rm -r storage/framework/cache/*
    # Alternativ (wenn php verfügbar): php artisan cache:clear
fi

# Log-Hinweis
bashio::log.info "Invoice Ninja stopped successfully. ✅"
