# rootfs/etc/s6-overlay/s6-rc.d/ninja/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Invoice Ninja
# Cleanup and graceful shutdown
# ==============================================================================

bashio::log.info "Invoice Ninja is stopping..."

# Versuche, PHP-Prozesse sauber zu beenden
if pgrep -f "php artisan serve" > /dev/null; then
    bashio::log.info "Stopping PHP server..."
    pkill -f "php artisan serve"
else
    bashio::log.info "No running PHP server found."
fi

# Cache und temporäre Dateien bereinigen (optional)
if [ -d /var/www/app/storage/framework/cache ]; then
    bashio::log.info "Clearing application cache..."
    rm -rf /var/www/app/storage/framework/cache/*
fi

# Log-Hinweis
bashio::log.info "Invoice Ninja stopped successfully. ✅"
bashio::log.info "All temporary files cleaned up."
