# rootfs/etc/s6-overlay/s6-rc.d/ninja/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Invoice Ninja
# Gracefully stop processes when the add-on is stopped
# ==============================================================================

bashio::log.info "Stopping Invoice Ninja..."

# Stop supervisord cleanly if it's running
if pgrep supervisord >/dev/null 2>&1; then
    bashio::log.info "Stopping supervisord..."
    pkill -SIGTERM supervisord
    sleep 5
fi

# Clean temporary files
bashio::log.info "Cleaning temporary files..."
rm -rf /tmp/* || true

bashio::log.info "Invoice Ninja has been stopped successfully."
exit 0
