# rootfs/etc/s6-overlay/s6-rc.d/ninja/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Invoice Ninja
# Starts Invoice Ninja under s6 supervision
# ==============================================================================

bashio::log.info "Setting up environment for Invoice Ninja..."

# Dynamically set APP_URL for ingress access
if bashio::addon.ingress_url; then
    export APP_URL=$(bashio::addon.ingress_url)
    bashio::log.info "APP_URL set to ${APP_URL}"
else
    export APP_URL="http://localhost"
    bashio::log.warning "No ingress URL found, using localhost."
fi

# Load addon configuration
export APP_KEY=$(bashio::config 'app_key')
export DB_HOST=$(bashio::services "mysql" "host")
export DB_PORT=$(bashio::services "mysql" "port")
export DB_DATABASE=$(bashio::config 'db_database')
export DB_USERNAME=$(bashio::config 'db_username')
export DB_PASSWORD=$(bashio::config 'db_password')

# Add custom host entry (optional)
echo "172.30.32.2 invoiceninja.test" >> /etc/hosts || bashio::log.warning "Could not write to /etc/hosts"

# Move to working directory
cd /var/www/app || bashio::exit.nok "Could not change to /var/www/app"

# Log configuration
bashio::log.info "Starting Invoice Ninja with:"
bashio::log.info "  APP_URL=${APP_URL}"
bashio::log.info "  DB_HOST=${DB_HOST}"
bashio::log.info "  DB_DATABASE=${DB_DATABASE}"
bashio::log.info "  DB_USERNAME=${DB_USERNAME}"

# Start supervisord through the official entrypoint
exec /usr/local/bin/docker-entrypoint supervisord -n
